# docker/resm_supervisor.conf
[program:resm]
command=/usr/local/bin/resm_shell_runner.sh
stderr_logfile = /var/log/supervisor/resm.log
stdout_logfile = /var/log/supervisor/resm.log
